# deploy/firebat-paneld.service
[Unit]
Description=Firebat Panel Daemon - LCD and LED control for Firebat S1 / AceMagic S1
Documentation=https://github.com/P4LMTR33/firebat-s1-panel
After=network.target dbus.service
Wants=dbus.service

[Service]
Type=simple
ExecStart=/usr/local/bin/firebat-paneld
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

# Logging
SyslogIdentifier=firebat-paneld
StandardOutput=journal
StandardError=journal

# State directory (creates /var/lib/firebat-panel)
StateDirectory=firebat-panel

# Security hardening
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/firebat-panel
PrivateTmp=yes
NoNewPrivileges=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes

# We need access to USB devices and /sys for sensors
DeviceAllow=/dev/ttyUSB* rw
DeviceAllow=/dev/hidraw* rw
SupplementaryGroups=plugdev dialout

[Install]
WantedBy=multi-user.target
