# deploy/firebat-paneld.service
[Unit]
Description=Firebat Panel Daemon - LCD and LED control for Firebat S1 / AceMagic S1
Documentation=https://github.com/P4LMTR33/firebat-s1-panel
After=network.target dbus.service
Wants=dbus.service
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/firebat-paneld
Restart=on-failure
RestartSec=5

# Logging
SyslogIdentifier=firebat-paneld
StandardOutput=journal
StandardError=journal

# State directory
StateDirectory=firebat-panel

# Minimal security to ensure hardware access
ProtectSystem=no
ProtectHome=no
PrivateTmp=no
# No restrictive device policies
DevicePolicy=auto

[Install]
WantedBy=multi-user.target
